# net-firewall/cs-firewall-bouncer-bin/files/crowdsec-firewall-bouncer.openrc
#!/sbin/openrc-run

command="crowdsec-firewall-bouncer"
command_args="-c ${CFG}crowdsec-firewall-bouncer.yaml"
pidfile=

name="crowdsec-firewall-bouncer"
description="The firewall bouncer for CrowdSec"
