# net-analyzer/cs-firewall-bouncer-bin/files/crowdsec.openrc
#!/sbin/openrc-run

command=/usr/local/bin/crowdsec
command_args="-c /etc/crowdsec/config.yaml"
pidfile=

name="crowdsec"
description="Crowdsec agent"
