# net-analyzer/cs-aws-waf-bouncer/files/cs-aws-waf-bouncer.service
[Unit]
Description=The AWS WAF bouncer for CrowdSec
After=syslog.target network.target remote-fs.target nss-lookup.target crowdsec.service
Before=netfilter-persistent.service

[Service]
Type=notify
ExecStart=/usr/bin/crowdsec-aws-waf-bouncer -c /etc/crowdsec/crowdsec-aws-waf-bouncer.yaml
ExecStartPre=/usr/bin/crowdsec-aws-waf-bouncer -c /etc/crowdsec/crowdsec-aws-waf-bouncer.yaml -t
ExecStartPost=/bin/sleep 0.1

[Install]
WantedBy=multi-user.target
