# net-analyzer/crowdsec/files/crowdsec.openrc
#!/sbin/openrc-run

command=/usr/bin/crowdsec
command_args="-c /etc/crowdsec/config.yaml"
#pidfile=
name="crowdsec"
description="Crowdsec agent"
command_background="yes"
pidfile="/run/$RC_SVCNAME/$RC_SVCNAME.pid"

start_pre () {
  checkpath --directory --owner $command_user:$command_user --mode 0755 \
  /run/$RC_SVCNAME
}
