# net-analyzer/crowdsec/files/crowdsec.openrc
#!/sbin/openrc-run

command=/usr/bin/crowdsec
command_args="-c /etc/crowdsec/config.yaml"
pidfile="/run/$RC_SVCNAME/$RC_SVCNAME.pid"
command_background="yes"

name="crowdsec"
description="Crowdsec agent"
