# net-analyzer/cs-blocklist-mirror/files/cs-blocklist-mirror.openrc
#!/sbin/openrc-run

command="/usr/bin/crowdsec-blocklist-mirror"
command_args="-c /etc/crowdsec/crowdsec-blocklist-mirror.yaml"
pidfile=
supervisor="supervise-daemon"

name="crowdsec-blocklist-mirror"
description="CrowdSec Blocklist Mirror"
